# gamescope_session/usr/lib/systemd/system/gamescope_steam_permissions.service
[Unit]
Description=change non-udev permissions for gamescope session

[Service]
Type=oneshot
ExecStart=/usr/bin/bash /usr/libexec/gamescope_steam_permissions.sh

